# deploy/cryptotrader.service
[Unit]
Description=Cryptocurrency Trading Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=cryptotrader
Group=cryptotrader
WorkingDirectory=/opt/cryptotrader
EnvironmentFile=-/opt/cryptotrader/.env
ExecStart=/opt/cryptotrader/venv/bin/python -m cryptotrader.main
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
